# storescptc.socket
[Unit]
Description = DICOM receiver

[Socket]
ListenStream = 11114
Accept = yes

[Install]
WantedBy = sockets.target
