# storescpfl.socket
[Unit]
Description = DICOM receiver

[Socket]
ListenStream = 11115
Accept = yes

[Install]
WantedBy = sockets.target
