# storescpfl.socket
[Unit]
Description = DICOM receiver

[Socket]
ListenStream = 11112
Accept = yes

[Install]
WantedBy = sockets.target
