# storescpfl@.service
[Unit]
Description=DIT3 DICOM receiver service

[Service]
ExecStart=storescp -xs -fe .dcm --inetd -od '/home/fuentes/wolfpacs/fl' --exec-on-reception '/home/fuentes/wolfpacs/storeandorganize.sh #p/#f fl'
StandardInput=socket
User=fuentes
